# launchd/com.openclaw.journal_export.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.journal_export</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/chenchinghung/openclaw/options/journal_export.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/chenchinghung/openclaw/options</string>

    <!-- 每月 1 日 06:00 跑（預設匯出上月）；events_sync 03:30、event_analysis 04:00 都已跑完 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/chenchinghung/openclaw/options/logs/journal_export.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/chenchinghung/openclaw/options/logs/journal_export.log</string>
</dict>
</plist>
